# systemd/claude-memory-orphan-backfill.timer
[Unit]
Description=Fire orphan-edge backfill 4× per day (00/06/12/18)
Documentation=https://github.com/vbcherepanov/total-agent-memory

[Timer]
# Equivalent to macOS StartCalendarInterval 00:00/06:00/12:00/18:00
OnCalendar=*-*-* 00,06,12,18:00:00
Persistent=true
Unit=claude-memory-orphan-backfill.service

[Install]
WantedBy=timers.target
